# motd/gitkeeper-motd.timer
[Unit]
Description=Generate git-keeper MOTD (timer)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
